# l4d2/rootfs/etc/s6-overlay/s6-rc.d/init-l4d2-config/up
/etc/s6-overlay/s6-rc.d/init-l4d2-config/run