# l4d2/rootfs/etc/s6-overlay/s6-rc.d/svc-l4d2/run
#!/usr/bin/with-contenv bash

/usr/bin/with-contenv bash /opt/gameserver/scripts/l4d2/run.sh
